# supervisor.ini
[supervisord]
nodaemon = true
user = root

[program: wiremock]
command = /bin/bash -c "exec /run-wiremock.sh"
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr = true

[program: envoy-proxy-template]
user = root
command = /bin/bash -c "exec /run-envoy.sh"
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr = true

